# setup_iptables.conf
start on startup
task
exec /opt/setup_iptables.sh
